# src/deb/control/postinst
#!/bin/sh

set -e

update-rc.d [[artifactId]]_[[version]] defaults
service [[artifactId]]_[[version]] start

exit 0