# src/deb/control/postrm
#!/bin/sh

set -e

update-rc.d -f [[artifactId]]_[[version]] remove

if [ -e /var/log/[[artifactId]]_[[version]].log ]; then
  rm /var/log/[[artifactId]]_[[version]].log
fi

exit 0